# mss-master.service
[Unit]
Description=MSS ROS2 Master Service
After=network.target
Wants=network.target

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/mss_ros/src
ExecStart=/opt/ros/humble/bin/ros2 launch mss_bringup all_nodes.launch.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
Environment=ROS_DOMAIN_ID=0

[Install]
WantedBy=multi-user.target
